# vminit/mon/mon/smf/grafana.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- 
    Manifest automatically generated by smfgen.
 -->
<service_bundle type="manifest" name="application-grafana" >
    <service name="application/grafana" type="service" version="1" >
        <create_default_instance enabled="true" />
        <dependency name="dep0" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>
	<exec_method type="method" name="start" exec="/mon/grafana/grafana/bin/grafana-server --homepath /mon/grafana/grafana --config /mon/grafana/grafana.ini  &amp;" timeout_seconds="60" >
            <method_context working_directory="/mon/grafana" >
                <method_credential user="mon" />
            </method_context>
        </exec_method>
        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60" >
            <method_context working_directory="/mon/grafana" >
                <method_credential user="mon" />
            </method_context>
        </exec_method>
        <template >
            <common_name >
                <loctext xml:lang="C" >Grafana</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
