# vminit/cockroachdb.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- 
    Manifest automatically generated by smfgen.
 -->
<service_bundle type="manifest" name="application-cockroachdb" >
    <service name="application/cockroachdb" type="service" version="1" >
        <create_default_instance enabled="false" />
        <dependency name="dep0" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>
        <dependency name="dep1" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/system/filesystem/local:default" />
        </dependency>
        <dependency name="dep2" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/network:default" />
        </dependency>
        <dependency name="dep3" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/name-services:default" />
        </dependency>
        <exec_method type="method" name="start" exec="/cockroachdb/bin/cockroach start --insecure --store /cockroachdb/data --listen-addr=%{config/my_internal_ip} --join=%{config/other_internal_ips} --background &amp;" timeout_seconds="60" >
            <method_context working_directory="/" >
                <method_credential user="cockroachdb" />
            </method_context>
        </exec_method>
        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" />
        <property_group name="config" type="application">
          <propval name="my_internal_ip" type="astring" value="" />
          <propval name="other_internal_ips" type="astring" value="" />
        </property_group>
        <template >
            <common_name >
                <loctext xml:lang="C" >CockroachDB node</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
