# vminit/ntpfix.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="application-ntpfix" >
    <service name="application/ntpfix" type="service" version="1" >
        <create_default_instance enabled="false" />

        <dependency name="dep0" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/multi-user:default" />
        </dependency>
        <dependency name="dep1" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/system/filesystem/local:default" />
        </dependency>
        <dependency name="dep2" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/network:default" />
        </dependency>
        <dependency name="dep3" grouping="require_all" restart_on="error" type="service" >
            <service_fmri value="svc:/milestone/name-services:default" />
        </dependency>

        <!-- At boot, we want to start before NTP. -->
        <dependent name="ntpfix-ntp" grouping="require_all" restart_on="none">
          <service_fmri value="svc:/network/ntp:default" />
        </dependent>

        <!--
          Unfortunately, it appears that ntpdate returns exit status 1 even when
          this works.  We'll exacerbate this dubious hack by ignoring the exit
          status.  Piping to `cat` is the easiest way to do this. (See
          https://www.illumos.org/issues/13076.)  Chrony can't come
          fast enough.
        -->
        <exec_method type="method" name="start" exec="ntpdate -b %{config/ntp_server} | cat" timeout_seconds="60" />
        <exec_method type="method" name="stop" exec=":true" timeout_seconds="10" />
        <property_group name="config" type="application">
          <propval name="ntp_server" type="astring" value="169.254.169.123" />
        </property_group>
        <property_group name="startd" type="framework">
          <propval name="duration" type="astring" value="transient" />
        </property_group>
        <template >
            <common_name >
                <loctext xml:lang="C" >NTP fix</loctext>
            </common_name>
        </template>
    </service>
</service_bundle>
